# stageApplication/00-application/files/shop-backup.service
[Unit]
Description=Shop Inventory Database Backup Service
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/shop-backup.sh
User=root
Group=root
Environment="DJANGO_BACKUP_PASSWORD=${DJANGO_BACKUP_PASSWORD}"

[Install]
WantedBy=multi-user.target
