# stageAccessPoint/00-accesspoint/files/shop-wifi-setup.service
[Unit]
Description=Shop WiFi Hotspot Setup
After=network.target NetworkManager.service dnsmasq.service
Wants=dnsmasq.service
Before=shop-inventory.service

[Service]
Type=oneshot
EnvironmentFile=/etc/shop-inventory/config
# Using nmcli dev wifi hotspot command for simplified hotspot creation
ExecStart=/bin/bash -c '\
    if [ -z "$WIFI_SSID" ] || [ -z "$WIFI_PASS" ]; then \
        echo "ERROR: WIFI_SSID and WIFI_PASS must be set in config"; \
        exit 1; \
    fi && \
    if ! nmcli -q c show "Hotspot"; then \
        nmcli dev wifi hotspot ifname wlan0 ssid "$WIFI_SSID" password "$WIFI_PASS"; \
    else \
        nmcli c up "Hotspot"; \
    fi && \
    sleep 2 && \
    systemctl restart dnsmasq'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
